# roles/osp-www/templates/supervisor/osp-www-dev.conf.j2

[program:osp-www-dev]
command={{ osp_src }}/env/bin/python {{ osp_src }}/osp/www/app.py
autostart=false
