# roles/osp-www/templates/supervisor/osp-www-uwsgi.conf.j2

[program:osp-www-uwsgi]
command={{ osp_src }}/env/bin/uwsgi --ini {{ uwsgi_config }}/osp-www.ini
autostart=false
