# roles/supervisor/templates/program.conf.j2

[program:{{ item.name }}]

command={{ item.command }}

autostart={{ item.autostart | default(false) }}

{% if item.numprocs is defined %}
numprocs={{ item.numprocs }}
{% endif %}
