# roles/supervisor/templates/program.conf.j2

[program:{{ item.name }}]

command={{ item.command }}

autostart={{ item.autostart | default(false) }}

numprocs={{ item.numprocs | default(1) }}
