# roles/osp-worker/templates/supervisor/osp-worker.conf.j2

[program:osp-worker]
command={{ osp_src }}/env/bin/uwsgi --ini {{ uwsgi_config }}/osp-worker.ini
autostart=false
