# roles/osp-worker/templates/supervisor/rq-worker.conf.j2

[program:rq-worker]
command={{ osp_src }}/env/bin/rqworker
process_name=%(program_name)s_%(process_num)02d
numprocs={{ osp_numprocs }}
autostart=false
