# roles/osp/templates/supervisor/tika-server.conf.j2

[program:tika-server]
command=java -jar {{ osp_tika_path }}/{{ osp_tika_jar }}
autostart=false
