# roles/osp-web/templates/supervisor/osp-web-uwsgi.conf.j2

[program:osp-web-uwsgi]
command={{ osp_src }}/env/bin/uwsgi --ini {{ uwsgi_config }}/osp-web.ini
autostart=false
